# deploy/aiaspeech.service
[Unit]
Description=AIASpeech Flask API
After=network.target mysql.service

[Service]
Type=notify
User=root
Group=www-data
WorkingDirectory=/var/www/aiaspeech/backend
Environment="PATH=/var/www/aiaspeech/venv/bin"
ExecStart=/var/www/aiaspeech/venv/bin/gunicorn \
    --config /var/www/aiaspeech/deploy/gunicorn_config.py \
    --workers 4 \
    --bind 127.0.0.1:5000 \
    run:app

Restart=always
RestartSec=10

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
